# governance/archive/tasks/project-history/legacy-project-scripts/scripts/setup/telegram-bot.service
[Unit]
# 示例 systemd 单元，请按部署环境调整 User / WorkingDirectory / EnvironmentFile
Description=FateCat Telegram Bot
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=fatecat
WorkingDirectory=/srv/fatecat
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=/etc/fatecat/fatecat.env
ExecStart=/srv/fatecat/.venv/bin/fatecat serve bot
Restart=always
RestartSec=5
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
